# dnssec-triggerd.service.in
[Unit]
Description=Reconfigure local DNSSEC resolver on connectivity changes
After=NetworkManager.service unbound.service dnssec-triggerd-keygen.service
Requires=unbound.service
Wants=dnssec-triggerd-keygen.service

[Service]
PIDFile=@pidfile@
Type=simple
Restart=always
ExecStart=@sbindir@/dnssec-triggerd -d
ExecStartPre=-@libexecdir@/dnssec-trigger-script --prepare
ExecStartPost=-@libexecdir@/dnssec-trigger-script --update
ExecStopPost=-@libexecdir@/dnssec-trigger-script --cleanup

[Install]
WantedBy=multi-user.target
