# osx/dnssec-trigger-panel.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
	"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>dnssec-trigger-panel</string>
	<key>ProgramArguments</key>
	<array>
		<string>0bindir0/dnssec-trigger-panel</string>
	</array>
</dict>
</plist>
