# osx/nl.nlnetlabs.dnssec-triggerd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC -//Apple Computer//DTD PLIST 1.0//EN http://www.apple.com/DTDs/PropertyList-1.0.dtd >
<plist version="1.0">
<!-- 
plist file for use by MacOSX launchd(8) using launchctl(1).
Copy this file to /Library/LaunchDaemons.  Launchd keeps the daemon running.
User root because it has to write to resolv.conf.
-->
<dict>
	<key>Label</key>
	<string>nl.nlnetlabs.dnssec-triggerd</string>

	<key>ProgramArguments</key>
	<array>
		<string>@sbindir@/dnssec-triggerd</string>
		<string>-d</string>
	</array>

	<key>UserName</key>
	<string>root</string>

	<key>RunAtLoad</key>
	<true/>

</dict>
</plist>
